# deploy/game_module.service
[Unit]
Description=Game Module API
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/project
Environment="GAME_MODULE_API_KEY=your-secret"
Environment="GM_MAX_REQUEST_SIZE=4194304"
ExecStart=/path/to/venv/bin/uvicorn game_module.http_api:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
